# ansible/files/ttyd-vtysh.target
[Unit]
Description=ttyd-vtysh target

[Install]
WantedBy=default.target
