# ansible/files/systemd/ttyd-container.target
[Unit]
Description=ttyd-container target

[Install]
WantedBy=default.target
